# Qt/GuidoServer/guidoweb.conf
# guidohttpd is a remote web service that provides the
# GUIDO Engine library API over the web
# this file is intended as a config file for upstart
#

description "Remote music sore web service based on the GUIDO Engine"
author      "Mike Solomon - Dominique Fober"


env PATH=/usr/local/bin:/usr/bin:
chdir /home/user/bin
setuid user

start on runlevel [2345]
stop  on runlevel [016]

exec /home/user/bin/guidohttpserver
